# sidekick/scripts/com.sidekick.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sidekick.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/elie/Projects/AIDevTools/sidekick/bin/sidekick</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/elie/Projects/AIDevTools/sidekick/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/elie/Projects/AIDevTools/sidekick/logs/stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/elie/Projects/AIDevTools/sidekick</string>
</dict>
</plist>